# database_exporting/schedule.crontab
0 0 * * 0 cd ~/npm-follower && ./database_exporting/export_metadata_external.sh
